# ejbca/test/config/supervisor/supervisord.conf
[supervisord]
# for test
nodaemon=false

# mysql              
[program:mysql]      
command=/bin/sh /usr/bin/mysql_up.sh