# ejbca/test/config/supervisor/supervisord.conf
[supervisord]
# for test
nodaemon=true

# mysql              
[program:mysql]      
command=/bin/sh /usr/bin/mysql_up.sh
priority = 10

# ejbca
#[program:ejbca]
#command=/bin/sh /usr/bin/build-ejbca.sh