# mysql/supervisord.conf
[supervisord]
nodaemon=true

# mysql              
[program:mysql]      
command=/bin/sh /mysql/mysql_up.sh